# home/config/systemd/user/clean-sessions.service
[Unit]
Description=Clean old session files regularly

[Service]
ExecStart=/home/glen/bin/clean-sessions
WorkingDirectory=/home/glen
