# home/config/systemd/user/fullscreen-anti-timeout.service
[Unit]
Description=Prevent lockscreen timeouts when in fullscreen applications

[Service]
ExecStart=/home/glen/bin/xautolock-wiggle-if-fullscreen
